# nova-v3.service
[Unit]
Description=NOVA v3 - AI Agent System
Documentation=https://github.com/NVA91/b-nova-v3
# Wichtig: Erst starten, wenn Docker wirklich bereit ist
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Wir nutzen eine Variable für den Pfad (definiert in Ansible)
WorkingDirectory={{ app_install_dir | default('/opt/b-nova-v3') }}

# V2 Syntax (docker compose) statt V1 (docker-compose)
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300

# Restart policy: Versuche es neu, wenn es crasht
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target